# Scripts/crontab
# TO USE THIS CRONTAB RUN:
# sudo crontab -u nick crontab

*/10 * * * * sh /Users/nick/.config/nix/Scripts/rclone_sync.bash >& "/Users/nick/Logs/rsync_bisync_mount/$(date)-rsync-bisync.txt"
