# Scripts/crontab
*/10 * * * * /Users/nick/.config/nix/Scripts/rsync_bisync.sh
